# printer-screensaver.conf
[program:np_screensaver]
command=/home/pi/.local/bin/uv run screensaver.py
directory=/home/pi/pi-thermal
user=pi
