# printer-spool.conf
[program:np_printer]
command=/home/pi/.local/bin/uv run main.py
directory=/home/pi/pi-thermal
user=pi
