# printer-spool.conf
[program:np_printer]
command=/home/pi/.local/bin/uv run uvicorn main:app --host 0.0.0.0 --port 4201
directory=/home/pi/pi-thermal
user=pi
